# src-config/files/etc/uci-defaults/01-init-octoprint-permissions
#!/bin/sh

mkdir -p /opt/octoprint
chown -R octop:octop /opt/octoprint
chmod -R ug=rwx /opt/octoprint

chown -R :octop /etc/octoprint
chmod -R g=rw /etc/octoprint

